# rootfs/etc/cron.d/peering-manager
30 2 * * * root cd /opt/peering-manager && python3 manage.py peeringdb_sync
55 * * * * root cd /opt/peering-manager && python3 manage.py configure_routers
0  * * * * root cd /opt/peering-manager && python3 manage.py poll_peering_sessions --all
0  0 * * * root cd /opt/peering-manager && python3 manage.py check_for_ix_peering_sessions
30 4 * * * root cd /opt/peering-manager && python3 manage.py grab_prefixes
